# target/linux/ipq40xx/base-files/etc/init.d/cpufreq
#!/bin/sh /etc/rc.common

START=15

target_min=500000

get_min_freq() {
	local available=$(cat /sys/devices/system/cpu/cpufreq/policy0/scaling_available_frequencies)
	for freq in $available; do
		[ "$freq" -lt "$target_min" ] && continue
		echo $freq
		return
	done
}

boot() {
	local min_freq

	min_freq=$(get_min_freq)

	if [ -z "$min_freq" ]; then
		logger -t cpufreq "No frequency above $((target_min/1000))MHz found, leaving default"
		return
	fi

	logger -t cpufreq "Setting min frequency to $((min_freq/1000))MHz"

	for policy in /sys/devices/system/cpu/cpufreq/policy*; do
		echo "$min_freq" > $policy/scaling_min_freq
	done

	if [ -f /sys/devices/system/cpu/cpufreq/ondemand/sampling_down_factor ]; then
		echo 10 > /sys/devices/system/cpu/cpufreq/ondemand/sampling_down_factor
		echo 50 > /sys/devices/system/cpu/cpufreq/ondemand/up_threshold
	fi
}
